# services/status-command-client.service
[Unit]
Description=Tavla command client
After=network-online.target
Wants=network-online.target

[Service]
User=pi
WorkingDirectory=/home/pi/tavla
ExecStart=/usr/bin/python3 /home/pi/tavla/command_client.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
